# supervisord.conf
[supervisord]
nodaemon = true

[program:crond]
command = /usr/sbin/crond -f

[program:polemarch]
autorestart = unexpected
startsecs = 0
command = sh -c "/opt/polemarch/bin/polemarchctl webserver"
